# deployment/secure-server-watcher.service.example
[Unit]
Description=Secure Server Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/secure-server-watcher
ExecStart=/opt/secure-server-watcher/run.sh --loop
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
